# test/test-execute/exec-monitor.service
[Unit]
Description=Test for ExecMonitor

[Service]
ExecStart=/bin/true
ExecMonitor=/bin/true
ExecMonitor=/bin/true
