# test/test-execute/exec-monitor-multiple.service
[Unit]
Description=Test for multiple ExecMonitor

[Service]
ExecStart=/bin/true
ExecMonitor=/bin/true
